# systemd/localnews_boot.timer
[Unit]
Description=Kicks off the Periodic DVR Timer

[Timer]
OnBootSec=1min
Unit=localnews_capture.timer

