# systemd/localnews_capture.service
[Unit]
Description=Local News Video Capture Script

[Service]
Type=simple
ExecStart=/opt/LocalNewsDVR/scripts/capture
RuntimeMaxSec=7300

