# systemd/localnews_boot.timer
[Unit]
Description=Prevents the DVR scheduler from starting within 5 minutes of bootup

[Timer]
OnBootSec=5min
Unit=localnews_capture.timer

