# systemd/localnews_hls.service
[Unit]
Description=Local News Video Packager Script

[Service]
Type=simple
ExecStart=/opt/LocalNewsDVR/scripts/hls


